# debian/linuxmuster-client-profile/etc/xdg/autostart/sshplus.py.desktop

[Desktop Entry]
Type=Application
Exec=/usr/bin/sshplus.py
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name[de_DE]=SSHPlus User-Menu
Name=SSHPlus User-Menu
Comment[de_DE]=Benutzermen端 f端r linuxmuster-client
Comment=Benutzermen端 f端r linuxmuster-client
